# files/FreeBSD/_ROOT_/usr/local/etc/periodic/daily/prune-be
#!/bin/sh

. ~/.config/walterjwhite/system-maintenance
timeout $_CONF_SYSTEM_MAINTENANCE_PRUNE_BE_TIMEOUT prune-be
