# files/FreeBSD/_ROOT_/usr/local/etc/periodic/daily/update-system
#!/bin/sh

. ~/.config/walterjwhite/system-maintenance
timeout $_CONF_SYSTEM_MAINTENANCE_UPDATE_SYSTEM_TIMEOUT update-system
