# Apple/files/_ROOT_/private/etc/periodic/daily/cleanup-auto-update-packages
#!/bin/sh

rm -f Downloads/Microsoft_AutoUpdate*.pkg
