# feature/zfs.feature/FreeBSD/files/_ROOT_/usr/local/etc/periodic/daily/zap
#!/bin/sh

zfs-zap-maintenance daily
zfs-rclone
